# cron/n8n_update.cron
# Exécuter la mise à jour de n8n tous les dimanches à 3h du matin
0 3 * * 0 ${ADMIN_DIR}/scripts/update_n8n.sh >> ${ADMIN_DIR}/logs/n8n_update.log 2>&1
