# cron/n8n_backup_db.cron
# Backup quotidien de la base de données n8n à 4h du matin
0 4 * * * ${ADMIN_DIR}/scripts/backup_n8n_db.sh >> ${ADMIN_DIR}/logs/n8n_backup_db.log 2>&1
